# agents/linux/network-optimizer-agent.service
[Unit]
Description=Network Optimizer Agent - System Metrics Collector
Documentation=https://github.com/yourusername/network-optimizer-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

# Path to the agent script
ExecStart=/opt/network-optimizer-agent/network-optimizer-agent.sh

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/network-optimizer-agent.log
StandardError=append:/var/log/network-optimizer-agent.log
SyslogIdentifier=network-optimizer-agent

# Environment variables (customize these)
Environment="INFLUXDB_URL=http://your-influxdb-host:8086"
Environment="INFLUXDB_TOKEN=your-token-here"
Environment="INFLUXDB_ORG=your-org"
Environment="INFLUXDB_BUCKET=network-metrics"
Environment="HOSTNAME=%H"
Environment="COLLECTION_INTERVAL=60"
Environment="COLLECT_SYSTEM_STATS=true"
Environment="COLLECT_DOCKER_STATS=false"
Environment="COLLECT_NETWORK_STATS=true"
Environment="NETWORK_INTERFACES=eth0"
Environment="DRY_RUN=false"
Environment="DEBUG=false"

# Security hardening
ProtectSystem=strict
ReadWritePaths=/var/log
PrivateTmp=true
NoNewPrivileges=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

[Install]
WantedBy=multi-user.target
